# deploy/supervisor_gunicorn.conf
; ═══════════════════════════════════════════════════════════════════
; Supervisor Config: Gunicorn (Flask App)
; Place at: /etc/supervisor/conf.d/idcard_gunicorn.conf
; ═══════════════════════════════════════════════════════════════════

[program:idcard_gunicorn]
; UPDATE the path below to match your Cloudways application directory
command=/home/master/applications/YOUR_APP_FOLDER/public_html/venv/bin/gunicorn app:app --bind 127.0.0.1:8000 --timeout 30 --workers 3 --threads 2 --access-logfile /home/master/applications/YOUR_APP_FOLDER/public_html/logs/gunicorn_access.log --error-logfile /home/master/applications/YOUR_APP_FOLDER/public_html/logs/gunicorn_error.log
directory=/home/master/applications/YOUR_APP_FOLDER/public_html
user=master
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
redirect_stderr=false
stderr_logfile=/home/master/applications/YOUR_APP_FOLDER/public_html/logs/gunicorn_err.log
stdout_logfile=/home/master/applications/YOUR_APP_FOLDER/public_html/logs/gunicorn_out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PATH="/home/master/applications/YOUR_APP_FOLDER/public_html/venv/bin:%(ENV_PATH)s"
